# com.jobalert.stock.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jobalert.stock</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/kimhakin/job-alert/.venv/bin/python3</string>
        <string>/Users/kimhakin/job-alert/stock_alert.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/kimhakin/job-alert/logs/stock.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kimhakin/job-alert/logs/stock-error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/kimhakin/job-alert</string>
</dict>
</plist>
